# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
loglevel=error

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:xvfb]
priority=10
command=/usr/bin/Xvfb :1 -screen 0 1440x1080x24
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/tmp/xvfb.log
redirect_stderr=true

[program:openbox]
priority=15
command=/usr/bin/openbox-session
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY=":1"
stdout_logfile=/tmp/openbox.log
redirect_stderr=true

[program:lxpanel]
priority=15
command=/usr/bin/lxpanel --profile LXDE
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY=":1"
stdout_logfile=/tmp/lxpanel.log
redirect_stderr=true

[program:pcmanfm]
priority=15
command=/usr/bin/pcmanfm --desktop --profile LXDE
autostart=true
autorestart=true
stopsignal=QUIT
environment=DISPLAY=":1"
stdout_logfile=/tmp/pcmanfm.log

[program:x11vnc]
priority=20
command=x11vnc -display :1 -xkb -forever -shared -repeat
autostart=true
autorestart=true
stopsignal=QUIT
stdout_logfile=/tmp/x11vnc.log
redirect_stderr=true

[group:vnc]
programs=xvfb,openbox,lxpanel,pcmanfm,x11vnc